# os-image/overlay/etc/systemd/system/kiacha-voice.service
[Unit]
Description=Kiacha Voice Runtime
After=kiacha-core.service network.target
Requires=kiacha-core.service

[Service]
Type=simple
ExecStart=/opt/kiacha/bin/kiacha-whisper --stream
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
